# scripts/jobhunt.service
[Unit]
Description=Remote SRE/Platform job scanner
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/remote-sre-job-scanner
EnvironmentFile=/opt/remote-sre-job-scanner/.env
ExecStart=/opt/remote-sre-job-scanner/.venv/bin/python jobhunt.py --config config.balanced.json scan

[Install]
WantedBy=multi-user.target
